# systemd/beagleg.service
[Unit]
Description=BeagleG - CNC motion control

[Service]
User=root
ExecStart=/usr/local/bin/run-beagleg.sh

[Install]
WantedBy=basic.target

